# golden_set.json
[
  {"question": "how do attackers pull passwords out of memory?", "expected": ["T1003.001"], "difficulty": "easy"},
  {"question": "what method is used when a malicious macro downloads a payload right after a user opens a phishing document?", "expected": ["T1204.002"], "difficulty": "easy"},
  {"question": "how do adversaries clear the windows event logs to hide their tracks?", "expected": ["T1685.005"], "difficulty": "easy"},
  {"question": "how can an operator execute commands by abusing the built-in windows management infrastructure?", "expected": ["T1047"], "difficulty": "easy"},
  {"question": "what technique involves compressing stolen files into a zip archive before stealing them?", "expected": ["T1560.001"], "difficulty": "easy"},
  {"question": "how do attackers use legitimate cloud storage apps like Dropbox or OneDrive to sneak data out of the network?", "expected": ["T1567.002"], "difficulty": "easy"},
  {"question": "how do attackers disable or uninstall endpoint security tools to avoid detection?", "expected": ["T1685"], "difficulty": "easy"},
  {"question": "what technique is used when an attacker injects malicious code into a normal running process like notepad?", "expected": ["T1055"], "difficulty": "medium"},
  {"question": "how do adversaries establish a footprint by creating a new account in the local Administrators group?", "expected": ["T1136.001"], "difficulty": "easy"},
  {"question": "how do attackers capture keystrokes to steal user credentials?", "expected": ["T1056.001"], "difficulty": "easy"},
  {"question": "how do threat actors perform a golden ticket attack?", "expected": ["T1558.001"], "difficulty": "hard"},
  {"question": "what technique lets an operator run malicious code by abusing the signed windows binary rundll32?", "expected": ["T1218.011"], "difficulty": "hard"},
  {"question": "how do adversaries trick a system into loading a malicious DLL by placing it next to the legitimate executable?", "expected": ["T1574.001"], "difficulty": "hard"},
  {"question": "how do attackers elevate privileges by bypassing windows user account control prompts?", "expected": ["T1548.002"], "difficulty": "medium"},
  {"question": "how do attackers find shared folders and drives available across a network?", "expected": ["T1135"], "difficulty": "easy"},
  {"question": "how do adversaries enumerate the user accounts that exist on a system or domain?", "expected": ["T1087"], "difficulty": "easy"},
  {"question": "how does malware figure out what operating system and hardware it is running on?", "expected": ["T1082"], "difficulty": "easy"},
  {"question": "how do attackers blend their command-and-control traffic in with normal web browsing traffic?", "expected": ["T1071.001"], "difficulty": "medium"},
  {"question": "how do adversaries encrypt their command-and-control communications to avoid network inspection?", "expected": ["T1573"], "difficulty": "easy"},
  {"question": "how do attackers gain initial access by getting a user to click a malicious link in an email?", "expected": ["T1566.002"], "difficulty": "easy"},
  {"question": "how do adversaries break in by exploiting a vulnerability in an internet-facing web server?", "expected": ["T1190"], "difficulty": "easy"},
  {"question": "how does ransomware lock up an organization's files to extort payment?", "expected": ["T1486"], "difficulty": "easy"},
  {"question": "how do attackers delete volume shadow copies so victims cannot restore their files?", "expected": ["T1490"], "difficulty": "medium"},
  {"question": "how do adversaries take screenshots of a victim's desktop?", "expected": ["T1113"], "difficulty": "easy"},
  {"question": "how do attackers gather files of interest from a compromised machine's local drives?", "expected": ["T1005"], "difficulty": "easy"},
  {"question": "how do attackers steal data by sending it out through the same channel they use for command-and-control?", "expected": ["T1041"], "difficulty": "medium"},
  {"question": "how do attackers authenticate as a user with just the captured password hash, without cracking it?", "expected": ["T1550.002"], "difficulty": "hard"},
  {"question": "how do adversaries pull password hashes by pretending to be a domain controller and requesting replication?", "expected": ["T1003.006"], "difficulty": "hard"},
  {"question": "how do attackers request service tickets in order to crack service account passwords offline?", "expected": ["T1558.003"], "difficulty": "hard"},
  {"question": "how do adversaries proxy execution of malicious code through a signed microsoft binary to bypass application control?", "expected": ["T1218"], "difficulty": "medium"},
  {"question": "how do attackers steal and impersonate another user's access token to act with their privileges?", "expected": ["T1134.001"], "difficulty": "hard"},
  {"question": "how do adversaries hide malware by making it look like a legitimate system file or program?", "expected": ["T1036"], "difficulty": "medium"},
  {"question": "how do attackers stop defenders from getting alerts by turning off windows event logging?", "expected": ["T1685.001"], "difficulty": "medium"},
  {"question": "what method lets malware survive a reboot by registering itself as a windows service?", "expected": ["T1543.003"], "difficulty": "medium"},
  {"question": "how do attackers move laterally through an environment?", "expected": ["T1021", "T1570", "T1550", "T1563"], "difficulty": "expected-fail"},
  {"question": "give me a list of all techniques that involve modifying the windows registry.", "expected": ["T1547.001", "T1112", "T1574.011", "T1037.001"], "difficulty": "expected-fail"},
  {"question": "what are the differences between abusing scheduled tasks versus system services for persistence?", "expected": ["T1053.005", "T1543.003"], "difficulty": "expected-fail"},
  {"question": "which execution techniques do not require the use of command-line interfaces?", "expected": ["T1204", "T1106"], "difficulty": "expected-fail"}
]